# argo
#!/sbin/openrc-run

command="/bin/argo-tunnel.sh"

depend()
{
    need net
}
