# argo
#!/sbin/openrc-run

command="/usr/share/argo-tunnel.sh"

depend()
{
    need net
}
