# supervisord.conf
[supervisord]
nodaemon=true

[program:Argo]
command=/usr/share/argo-tunnel.sh
