# supervisord.conf
[supervisord]
nodaemon=true

#[program:traefik]
#command=/traefik

[program:Argo]
command=/usr/share/argo-tunnel.sh
